# dist/yiidmtmnode.service.dist
[Unit]
Description=yiidmtmnode
After=network-online.target mysql.service
Requires=mysql.service

[Service]
User=www-data
PIDFile=/media/www/node.local.net/html/console/runtime/daemon/daemon.pid
Type=forking
KillMode=process
ExecStart=/media/www/node.local.net/html/vendor/inpassor/yii2-daemon/yiid start
ExecStop=/media/www/node.local.net/html/vendor/inpassor/yii2-daemon/yiid stop

[Install]
WantedBy=multi-user.target
